# images/locust/resources/supervisor/locust.conf
[program:locust]
command=python /opt/locust/pod/pod.py
