# images/wormhole/gossip/resources/supervisor/gossip.conf
[program:gossip]
command=python /opt/gossip/gossip.py
