# images/flask-sample/resources/supervisor/flask.conf
[program:flask]
command=python /opt/flask/pod/pod.py
stderr_logfile=/var/log/supervisor/flask.stderr.log