# images/wormhole/mailer/resources/supervisor/mailer.conf
[program:mailer]
command=python /opt/mailer/mailer.py
