# images/cleaner/resources/supervisor/cleaner.conf
[program:cleaner]
command=python /opt/cleaner/pod/pod.py
