# images/watcher/resources/supervisor/watcher.conf
[program:watcher]
command=python /opt/watcher/pod/pod.py
